# docker/etc/s6-overlay/s6-rc.d/plexripper/run
#!/usr/bin/with-contenv sh
exec s6-setuidgid abc /usr/share/dotnet/dotnet /app/PlexRipper.WebAPI.dll